# eqdkp-api/eqdkp-api.conf
[program:eqdkp-api]
directory=/home/ubuntu/ro/
command=/home/ubuntu/ro/eqdkp-api
redirect_stderr=true
autorestart=true
stdout_logfile=/home/ubuntu/ro/eqdkp-api.log
stdout_logfile_maxbytes=5MB
